# target/2024-10-17_06:04:09_debian12-with-kde_spare/.config/autostart/login_script.bash.desktop
[Desktop Entry]
Exec=bash ~/.local/bin/kde/login_script.bash
Icon=dialog-scripts
Name=login_script.bash
Path=
Type=Application
X-KDE-AutostartScript=true
